# ops/launchd/com.therapyops.acorn-preflight.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.therapyops.acorn-preflight</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__REPO_ROOT__/scripts/automation/run_daily_automation.py</string>
    <string>--mode</string>
    <string>preflight</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>ACORN_AUTOMATION_REPORT_TO</key>
    <string>__REPORT_TO__</string>
    <key>ACORN_DOCKER_READY_TIMEOUT_SEC</key>
    <string>__DOCKER_TIMEOUT__</string>
    <key>ACORN_AUTOMATION_LOG_DIR</key>
    <string>__HOME__/Library/Logs/therapy-ops-agent</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/therapy-ops-agent/launchd-preflight.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/therapy-ops-agent/launchd-preflight.err.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>45</integer></dict>
  </array>
</dict>
</plist>
